# curator4-2012372389.conf
[program:curator4-2012372389]
command=python3 /home/curator4/curator4-2012372389/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-2012372389/test.err.log
stdout_logfile=/home/curator4/curator4-2012372389/test.out.log